# deploy/mindroom-widget.service
[Unit]
Description=MindRoom Frontend
After=network.target
Wants=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=PROJECT_USER
Group=PROJECT_GROUP
WorkingDirectory=PROJECT_DIR
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=BACKEND_PORT=8001
Environment=FRONTEND_PORT=3003

# Use the existing run-ui.sh script (no auto-update)
ExecStart=PROJECT_DIR/run-ui.sh

# Logging
StandardOutput=append:/var/log/mindroom-frontend.log
StandardError=append:/var/log/mindroom-frontend.log

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=PROJECT_DIR /var/log /tmp

# Kill process group on stop (to handle background processes)
KillMode=process

[Install]
WantedBy=multi-user.target
