# deploy/mindroom-widget.service
[Unit]
Description=MindRoom Configuration Widget (Frontend + Backend)
After=network.target
Wants=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=basnijholt
Group=users
WorkingDirectory=/home/basnijholt/Work/mindroom-2/widget
Environment=PATH=/home/basnijholt/.nix-profile/bin:/run/current-system/sw/bin:/usr/bin:/bin
Environment=BACKEND_PORT=8001
Environment=FRONTEND_PORT=3003

# Ensure we're on the latest main branch and deps are up to date
ExecStartPre=/bin/sh -c 'cd /home/basnijholt/Work/mindroom-2 && git fetch origin && git reset --hard origin/main'
ExecStartPre=/bin/sh -c 'cd /home/basnijholt/Work/mindroom-2/widget/backend && uv sync'
ExecStartPre=/bin/sh -c 'cd /home/basnijholt/Work/mindroom-2/widget/frontend && pnpm install'

# Use the existing run.sh script
ExecStart=/home/basnijholt/Work/mindroom-2/widget/run.sh

# Logging
StandardOutput=append:/var/log/mindroom-widget.log
StandardError=append:/var/log/mindroom-widget.log

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/basnijholt/Work/mindroom-2 /var/log /tmp

# Kill process group on stop (to handle background processes)
KillMode=process

[Install]
WantedBy=multi-user.target
