# deploy/mindroom-backend.service
[Unit]
Description=MindRoom Multi-Agent System Backend
After=network.target
Wants=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=5
User=PROJECT_USER
Group=PROJECT_GROUP
WorkingDirectory=PROJECT_DIR
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Start the mindroom backend (no auto-update)
ExecStart=PROJECT_DIR/.venv/bin/python -m mindroom.cli run --log-level INFO --storage-path PROJECT_DIR/tmp

# Logging
StandardOutput=append:/var/log/mindroom-backend.log
StandardError=append:/var/log/mindroom-backend.log

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=PROJECT_DIR /var/log /tmp

[Install]
WantedBy=multi-user.target
