# deploy/mindroom-autoupdate.timer
[Unit]
Description=MindRoom Auto-Update Timer
Requires=mindroom-autoupdate.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
